# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

cd /opt/radarr || exit

exec \
	s6-setuidgid root mono --debug Radarr.exe \
	-nobrowser -data=/config
